# src/conf/cron.d/ss-server.cron
0 */2 * * * /usr/local/bin/supervisorctl -c /etc/supervisord.conf restart all;